# clear-ltechagent-state.timer
[Unit]
Description=Run clear-ltechagent-state on boot and timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
