# clear-ltechagent-state.timer
[Unit]
Description=Run clear-ltechagent-state on timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
